# trading-bot.service
[Unit]
Description=MetaLabel Trading Bot
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/metalabel-trader
Environment="PATH=/home/ubuntu/metalabel-trader/venv/bin:/usr/bin"
# Use venv python if it exists, otherwise use system python
ExecStart=/bin/bash -c 'if [ -d venv ]; then venv/bin/python live_trading_bot.py; else python3 live_trading_bot.py; fi'
Restart=always
RestartSec=30
StandardOutput=append:/home/ubuntu/metalabel-trader/logs/bot.log
StandardError=append:/home/ubuntu/metalabel-trader/logs/bot-error.log

# Restart policy
StartLimitBurst=5
StartLimitIntervalSec=600

[Install]
WantedBy=multi-user.target
